# deploy/rpi/fukaibox-backend.service
[Unit]
Description=FukaiBox Backend API
After=network.target postgresql.service

[Service]
Type=simple
User=fukaibox
Group=fukaibox
WorkingDirectory=/opt/fukaibox/backend
Environment=PATH=/opt/fukaibox/venv/bin:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=-/opt/fukaibox/.env
ExecStart=/opt/fukaibox/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
